# touchscreen-autodisabler.service.in
[Unit]
Description=Touchscreen autodisabler

[Service]
Type=idle
BusName=wsid.Tad
ExecStartPre=@SLEEP@ 2
ExecStart=@bindir@/touchscreen-autodisabler

[Install]
WantedBy=default.target
